# bin/deploy/.crontab
# corn default crontab
# sec min hour day month weekday command
*/10 * * * * * echo hello-crontab
